# rpi/salt/axidraw/etc/profile.d/99-axidraw-aliases.sh
alias axicli='/usr/local/bin/axicli --port `/usr/bin/readlink -f /dev/$USER` --port_config 0'
