# rpi/salt/axidraw/etc/profile.d/99-axidraw-aliases.sh
# 99-axidraw-aliases
alias axicli='/opt/venv-axidraw/bin/axicli --port `/usr/bin/readlink -f /dev/$USER` --port_config 0'
